# webserver.init
#!/bin/sh /etc/rc.common

START=98
USE_PROCD=1

start_service() {
    procd_open_instance
    #procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param command /usr/bin/python3 /usr/bin/webserver/main.py
    procd_set_param respawn  # enables automatic restart on crash
    procd_set_param respawn_delay 5  # wait 5 seconds before restart
    procd_close_instance
}